# src/DEBIAN/postinst
chmod 644 /etc/systemd/system/rpi-fan-control.service
systemctl daemon-reload > /dev/null 2>&1 || true
systemctl enable rpi-fan-control.service
systemctl start rpi-fan-control.service
